# scripts/helloworld.service
[Unit]
description=Hello world Deamon

[Service]
Restart=on-failure
ExecStart=/usr/bin/node /home/ec2-user/helloworld.js
